# systemd/capacinator.service
[Unit]
Description=Capacinator Dev Server
Documentation=https://github.com/tony/Capacinator
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/capacinator
ExecStart=/usr/bin/node dist/server/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3456
Environment=DATABASE_PATH=/var/www/capacinator/data/capacinator.db
Environment=CORS_ORIGIN=http://dev.capacinator.com
Environment=LOG_LEVEL=info
Environment=TRUST_PROXY=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=capacinator

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/www/capacinator/data
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G

[Install]
WantedBy=multi-user.target